# supervisor.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile=/var/log/supervisor/nginx.out.log
priority=999

[program:gemini-api]
command=node /app/dist-server/server/index.js
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/gemini-api.err.log
stdout_logfile=/var/log/supervisor/gemini-api.out.log
environment=GEMINI_API_KEY="%(ENV_GEMINI_API_KEY)s",NODE_ENV="production",PORT="3000"
priority=998
